# home/.config/systemd/user/psd-resync.timer.d/frequency.conf
[Unit]
Description=Timer for Profile-sync-daemon - 5min

[Timer]
OnUnitActiveSec=
OnUnitActiveSec=5min
